# com.obsidian.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.obsidian.sync</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-c</string>
      <string>~/obsidian-sync.sh</string>
      <string>YOUR SCRIPT LOCATION</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>  <!-- Runs every 5 minutes -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/obsidian-sync.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/obsidian-sync.log</string>
  </dict>
</plist>

